# aviasales-bot.service
[Unit]
Description=Aviasales Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/TG_Aviasales
Environment=PATH=/home/ubuntu/TG_Aviasales/venv/bin
ExecStart=/home/ubuntu/TG_Aviasales/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 